# systemd/vibehost-demo.service
[Unit]
Description=Vibehost demo service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vibehost-demo
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
